# deployment/systemd/encrypted-backup.timer
[Unit]
Description=Run Encrypted Backup Daily at 03:00 UTC
Requires=encrypted-backup.service

[Timer]
OnCalendar=daily
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=encrypted-backup.service

[Install]
WantedBy=timers.target
